# system_files/usr/etc/profile.d/10-luka-aliases.sh
alias ls='eza'
alias ll='eza -lah'
alias cat='bat'
alias grep='rg'
